# config/includes.chroot/etc/xdg/autostart/notification-daemon.desktop
[Desktop Entry]
Type=Application
Name=Notification Daemon
Exec=/usr/lib/notification-daemon/notification-daemon
Terminal=false
NoDisplay=true
